# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
loglevel=info

[inet_http_server]
port=0.0.0.0:9001

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=http://localhost:9001

[program:01_ros2_launch]
command=/bin/bash -c "source /opt/ros/jazzy/setup.bash && source /home/ros2ws/install/setup.bash && ros2 launch mocap4r2_optitrack_driver optitrack2.launch.py"
autostart=false
autorestart=false
stdout_logfile=/var/log/supervisor/ros2_launch.log
stderr_logfile=/var/log/supervisor/ros2_launch.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

[program:02_lifecycle_activate]
command=/bin/bash -c "source /opt/ros/jazzy/setup.bash && source /home/ros2ws/install/setup.bash && ros2 lifecycle set /mocap4r2_optitrack_driver_node activate"
autostart=false
autorestart=false
startretries=0
stdout_logfile=/var/log/supervisor/lifecycle_activate.log
stderr_logfile=/var/log/supervisor/lifecycle_activate.log

[program:03_lifecycle_deactivate]
command=/bin/bash -c "source /opt/ros/jazzy/setup.bash && source /home/ros2ws/install/setup.bash && ros2 lifecycle set /mocap4r2_optitrack_driver_node deactivate"
autostart=false
autorestart=false
startretries=0
stdout_logfile=/var/log/supervisor/lifecycle_deactivate.log
stderr_logfile=/var/log/supervisor/lifecycle_deactivate.log

[program:04_lifecycle_cleanup]
command=/bin/bash -c "source /opt/ros/jazzy/setup.bash && source /home/ros2ws/install/setup.bash && ros2 lifecycle set /mocap4r2_optitrack_driver_node cleanup"
autostart=false
autorestart=false
startretries=0
stdout_logfile=/var/log/supervisor/lifecycle_cleanup.log
stderr_logfile=/var/log/supervisor/lifecycle_cleanup.log

[program:05_lifecycle_shutdown]
command=/bin/bash -c "source /opt/ros/jazzy/setup.bash && source /home/ros2ws/install/setup.bash && ros2 lifecycle set /mocap4r2_optitrack_driver_node shutdown"
autostart=false
autorestart=false
startretries=0
stdout_logfile=/var/log/supervisor/lifecycle_shutdown.log
stderr_logfile=/var/log/supervisor/lifecycle_shutdown.log
